# bing/bingscript.conf
start on runlevel [2345]
stop on runlevel [!2345]

exec /home/sud/Desktop/bing/bing.py